# app/supervisord.conf
[supervisord]
nodaemon=true  #在前台运行 supervisord，这在 Docker 容器中是必需的

# FastAPI 主应用服务
[program:fastapi]
command=uvicorn main:app --host 0.0.0.0 --port 8000  # 启动命令
directory=/app  # 命令执行的工作目录
user=root  # 运行用户
stdout_logfile=/app/logs/fastapi.log  # 标准输出日志文件
stderr_logfile=/app/logs/fastapi.err  # 标准错误日志文件
autostart=true  # supervisord 启动时自动启动该程序
autorestart=true  # 程序异常退出时自动重启

# 天气查询 MCP 服务
[program:weather_service]
command=uvicorn mcp_server.weather_service:app --host 0.0.0.0 --port 9001
directory=/app
user=root
stdout_logfile=/app/logs/weather_service.log
stderr_logfile=/app/logs/weather_service.err
autostart=true
autorestart=true

# 订单查询 MCP 服务
[program:order_service]
command=uvicorn mcp_server.order_service:app --host 0.0.0.0 --port 9002
directory=/app
user=root
stdout_logfile=/app/logs/order_service.log
stderr_logfile=/app/logs/order_service.err
autostart=true
autorestart=true
